# infra/cdk.out/asset.260b0498a747c39d086d5034f51cb9992c8a1e923b95b918b2e78a78375dd591/systemd/hamaro-backup.service
[Unit]
Description=Hamaro periodic live backup
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/hamaro/backup.sh
