# infra/cdk.out/asset.a780c760a3af3bda4a02fa660dfbd002185ef46dc4ab5ef6db370f489930429f/systemd/hamaro-watchdog.timer
[Unit]
Description=Run the Hamaro watchdog (fast heartbeat, ~1min RCON enrichment)

[Timer]
# boot.sh already fires an immediate heartbeat right after starting the
# container; this is just the fallback net in case that ever fails to run.
# 20s tick is cheap (Server List Ping, no RCON) — the RCON-heavy parts inside
# watchdog.sh self-time to ~60s regardless of how often this fires.
OnBootSec=20s
OnUnitActiveSec=20s
AccuracySec=3s

[Install]
WantedBy=timers.target
