# infra/cdk.out/asset.234cd2e7f97aed7e5ef2a32c24f76dceae26079f7389c38c7a5e9c07427b95ff/systemd/hamaro-boot.service
[Unit]
Description=Hamaro Minecraft boot (DNS update, config sync, server start)
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/usr/local/bin/hamaro-sync
ExecStart=/opt/hamaro/boot.sh
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
