# infra/cdk.out/asset.89cc4982fde848bcf7896a8dfd99d3ed5775521746a0d330de0a46e89e196e93/systemd/hamaro-watchdog.timer
[Unit]
Description=Run the Hamaro idle watchdog every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s

[Install]
WantedBy=timers.target
