# infra/cdk.out/asset.ae472ca895a9ddb620dc24a88898db03f53eb3c0ab113ac0f265e1655a7057a9/systemd/hamaro-watchdog.timer
[Unit]
Description=Run the Hamaro idle watchdog every minute

[Timer]
# boot.sh already fires an immediate heartbeat right after starting the
# container; this is just the fallback net in case that ever fails to run.
OnBootSec=30s
OnUnitActiveSec=1min
AccuracySec=15s

[Install]
WantedBy=timers.target
