# infra/cdk.out/asset.ae472ca895a9ddb620dc24a88898db03f53eb3c0ab113ac0f265e1655a7057a9/systemd/hamaro-watchdog.service
[Unit]
Description=Hamaro idle watchdog (one check)
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/hamaro/watchdog.sh
