# infra/cdk.out/asset.ae472ca895a9ddb620dc24a88898db03f53eb3c0ab113ac0f265e1655a7057a9/systemd/hamaro-backup.timer
[Unit]
Description=Back up the live world every 2 hours while the server is up

[Timer]
OnBootSec=2h
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
